# evaluation/properties/network_tc_setup.service
[Unit]
Description=Network Traffic Control Setup

[Service]
Type=oneshot
ExecStart=/etc/network_tc_setup.sh

[Install]
WantedBy=multi-user.target